# systemd/caper.service.in
[Unit]
Description=IPC-controlled wallpaper/widget utility for Hyprland
Documentation=https://github.com/Cameron-Cooper22/caper
PartOf=graphical-session.target
Requires=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/hyprpaper
Slice=session.slice
Restart=on-failure

[Install]
WantedBy=graphical-session.target
